# deploy/systemd/voice-bridge.service
[Unit]
Description=Voice Bridge FastAPI Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/voice-bridge
EnvironmentFile=/etc/voice-bridge.env
EnvironmentFile=-/etc/voice-bridge.secrets
ExecStart=/opt/voice-bridge/.venv/bin/uvicorn app.main:app --host ${VOICE_BRIDGE_HOST} --port ${VOICE_BRIDGE_PORT}
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
